# systemd/triage-tick.service.d/dispatch.conf
[Service]
# Production default: scharf. Der Unit-File selbst bleibt auf 0, damit ein
# frisch geklonter Host nicht versehentlich sofort dispatcht — dieser
# versionierte Drop-in flippt es für unseren Triage-Host ein.
Environment=TRIAGE_ENABLE_DISPATCH=1
# systemd setzt HOME nicht automatisch für root-Services. Ohne HOME findet
# gh /root/.config/gh/hosts.yml nicht und fällt auf "gh auth login" zurück
# → Engineer/Review-Dispatcher sterben mit NOPERMISSION bevor codex läuft.
Environment=HOME=/root
