# systemd/triage-tick.service.d/dispatch.conf
[Service]
# Production default: armed. The unit file itself stays at 0 so a
# newly cloned host doesn't accidentally start dispatching immediately — this
# versioned drop-in flips it to 1 for our production triage host.
Environment=TRIAGE_ENABLE_DISPATCH=1
# systemd doesn't automatically set HOME for root services. Without HOME,
# gh cannot find /root/.config/gh/hosts.yml and falls back to "gh auth login"
# -> Engineer/Review dispatchers fail with NOPERMISSION before codex runs.
Environment=HOME=/root
# Rendered by install.sh with the production host roots. Override
# TRIAGE_REPOS_DIR/TRIAGE_WORKTREES_DIR when installing on another host.
Environment="TRIAGE_REPOS_DIR=@TRIAGE_REPOS_DIR@"
Environment="TRIAGE_WORKTREES_DIR=@TRIAGE_WORKTREES_DIR@"
